# src/Easyjet/UmanBundle/install/uman/crontab
0 4 * * *	apache	/usr/bin/php /usr/share/swissinside/bin/console easyjet:sfsu-remove-expired-files
